# systemd/do-dyn-domain.timer
[Unit]
Description=Dynamically update Digital Ocean domain

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=multi-user.target
